# ai-assistant/configs/supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.stdout.log
stderr_logfile=/var/log/supervisor/nginx.stderr.log

[program:ai-assistant]
command=/usr/bin/env python /app/main.py
directory=/app
autorestart=true
stdout_logfile=/var/log/supervisor/ai-assistant.stdout.log
stderr_logfile=/var/log/supervisor/ai-assistant.stderr.log
